# server.conf
[program:server]
user=vagrant
autostart=true
autorestart=true
startretries=3
command=python server.py
directory=/vagrant
stdout_logfile=/vagrant/server.log
stderr_logfile=/vagrant/server.err
